# deployment/argos-v3.service
[Unit]
Description=Argos v3 (NVIDIA UI) Production Server on :5175
After=network.target

[Service]
Type=simple
User=__SETUP_USER__
Group=__SETUP_USER__
WorkingDirectory=__PROJECT_DIR__
# NOTE: EnvironmentFile= intentionally omitted. systemd's parser does NOT strip
# inline '# comment' from values, breaking Zod URL/number validation. The app's
# src/lib/server/env.ts calls dotenv.config() at module-import time, which DOES
# strip comments correctly. Mirrors argos-final.service.d/10-env.conf.
Environment="NODE_ENV=production"
Environment="NODE_OPTIONS=--max-old-space-size=1024"
Environment="PORT=5175"
Environment="HOST=0.0.0.0"
# Pre-spawn tmux-0 on startup so the v3 dashboard Terminal panel opens
# instantly (same as argos-final / argos-dev do for v1 / v2).
Environment="ARGOS_TERMINAL_PRESPAWN=1"
# Uses the same prod-server.ts entry as argos-final / argos-dev so the
# WebSocket upgrade for /terminal-ws is wired in v3 production. The
# src/hooks.server.ts port branch routes PORT=5175 to the /dashboard/v3
# (NVIDIA UI) surface. See memory feedback_port_ui_split_nonnegotiable.md.
ExecStart=__NODE_BIN__ --import tsx __PROJECT_DIR__/scripts/ops/prod-server.ts
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=argos-v3

[Install]
WantedBy=multi-user.target
