# deployment/argos-newui-dev.service
[Unit]
Description=Argos New UI Dev Server (active session-N branch) on :5175
After=network.target

[Service]
Type=simple
User=__SETUP_USER__
Group=__SETUP_USER__
Environment="PATH=/home/__SETUP_USER__/.nvm/versions/node/v22.22.2/bin:/usr/local/bin:/usr/bin:/bin"
# Worktree carrying the active v2 dev work. Daily flow: code on session-2,
# PR session-2 → dev, merge, refresh. The worktree-refresh@.timer keeps this
# branch synced if the user opts in (template instance points at this path).
WorkingDirectory=__PROJECT_DIR__-session-2
Environment="NODE_ENV=development"
Environment="NODE_OPTIONS=--max-old-space-size=2048"
# VITE_PORT is derived by start-vite.sh from $(basename $WorkingDirectory),
# so this defaults to :5175 when WorkingDirectory ends in `Argos-session-2`.
ExecStartPre=/bin/sh -c 'cd __PROJECT_DIR__-session-2 && git fetch origin --quiet || true'
ExecStart=__PROJECT_DIR__-session-2/scripts/dev/start-vite.sh
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=argos-newui-dev

[Install]
WantedBy=multi-user.target
