# deployment/argos-worktree-refresh@.service
[Unit]
Description=Argos worktree refresh — fast-forward %i to origin if clean
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=jetson2
Group=jetson2
Environment="PATH=/home/jetson2/.nvm/versions/node/v22.22.2/bin:/usr/local/bin:/usr/bin:/bin"
# %i is the unit instance — pass it to the script as the worktree path.
# Escape slashes in instance name with `-` per systemd-escape(1):
#   sudo systemctl enable argos-worktree-refresh@$(systemd-escape -p /home/jetson2/code/Argos).timer
ExecStart=/home/jetson2/code/Argos/scripts/ops/worktree-refresh.sh %I
StandardOutput=journal
StandardError=journal
SyslogIdentifier=argos-worktree-refresh

[Install]
WantedBy=multi-user.target
